# chatgpt-wp-bridge.service
[Unit]
Description=WordPress Publishing Bridge
After=network.target

[Service]
Type=simple
User=wpbridge
Group=wpbridge
WorkingDirectory=/opt/wp-bridge
EnvironmentFile=/opt/wp-bridge/.env
ExecStart=/opt/wp-bridge/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=5
StandardOutput=append:/var/log/wp-bridge.log
StandardError=append:/var/log/wp-bridge.log

[Install]
WantedBy=multi-user.target
