# deploy/chatgpt-wp-bridge.service
[Unit]
Description=ChatGPT WordPress Publishing Bridge
After=network.target

[Service]
Type=simple
User=chatgptwpbridge
Group=chatgptwpbridge
WorkingDirectory=/opt/chatgpt-wp-bridge
EnvironmentFile=/opt/chatgpt-wp-bridge/.env
ExecStart=/opt/chatgpt-wp-bridge/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=5
StandardOutput=append:/var/log/chatgpt-wp-bridge.log
StandardError=append:/var/log/chatgpt-wp-bridge.log

[Install]
WantedBy=multi-user.target
